# meta-atm/recipes-hipos/atm/files/atm.service
[Unit]
Description=Atm

[Service]
ExecStart=/usr/bin/atm
Restart=always

[Install]
WantedBy=multi-user.target