# meta-atm/recipes-hipos/atm/files/atm.service
[Unit]
Description=Atm

[Service]
ExecStart=/etc/atm/start.sh
Restart=always

[Install]
WantedBy=multi-user.target